# hbase-hbck.timer
[Unit]
Description=HBase hbck timer

[Timer]
OnCalendar=*-*-* *:20:00
Unit=hbase-hbck.service

[Install]
WantedBy=multi-user.target
