# hbase-hbck.service
[Unit]
User=hbase
Description=HBase hbck

[Service]
Type=Simple
ExecStart=/usr/bin/hbase hbck

[Install]
WantedBy=multi-user.target
