# hbase-hbck.service
[Unit]
Description=HBase hbck

[Service]
User=hbase
Type=simple
ExecStart=/usr/bin/hbase hbck

[Install]
WantedBy=multi-user.target
